# tsoai-deploy.service
[Unit]
Description=Auto Deploy TSOAI App

[Service]
Type=oneshot
ExecStart=/opt/tsoai/deploy.sh
